# etc/cron.daily/420.status-network
#!/bin/sh
#
# $FreeBSD$
#


echo "Network interface status:"

flags="${daily_status_network_netstat_flags}"

ip -stats ${flags} link && rc=0 || rc=3

echo ""

exit $rc
