# etc/cron.daily/450.status-security
#!/bin/sh
#
# $FreeBSD$
#

run-parts /etc/cron.security && rc=1 || rc=3

exit $rc
